# monitor/dry-age-monitor.service
[Unit]
Description=BME690 dry-aging fridge logger
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=tan
SupplementaryGroups=i2c
WorkingDirectory=/home/tan/dry-age-monitor
ExecStart=/home/tan/dry-age-monitor/monitor/.venv/bin/python /home/tan/dry-age-monitor/monitor/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
